# dist/Network/template_ookla_speedtest_by_zabbix_agent/7.2/files/zabbix-speedtest.timer
[Unit]
Description=Run Ookla speedtest every 2h

[Timer]
OnBootSec=5min
OnUnitActiveSec=2h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
