# dist/Network/template_ookla_speedtest_by_zabbix_agent/7.2/files/zabbix-speedtest.service
[Unit]
Description=Refresh Ookla speedtest cache
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=zabbix
StateDirectory=zabbix-speedtest
EnvironmentFile=-/etc/default/zabbix-speedtest
ExecStart=/usr/local/bin/zbx-speedtest.sh run
TimeoutStartSec=5min
